# debian/postrm
#!/bin/sh

set -e

update-rc.d llmbs remove >/dev/null

exit 0