# debian/postinst
#!/bin/sh

set -e 

update-rc.d llmbs defaults >/dev/null
# source /etc/profile.d/llmbs.sh
service llmbs start

exit 0
