# debian/postinst
#!/bin/sh

set -e 

update-rc.d llmbs defaults >/dev/null


exit 0